# sys/wod-npm-app.service.j2
[Unit]
Description=Launching WoD npm based service
After=syslog.target network.target

[Service]
User={{ WODUSER }}
ExecStart=/usr/local/bin/wod-npm-app

[Install]
WantedBy=multi-user.target
